# overlay/system/etc/init/wlan-start.conf
start on started lightdm

script

echo ON > /dev/wlan

end script
